# munin/files/crontab
*/5 * * * * /usr/local/bin/munin-cron > /dev/null
