# lxc-solr/templates/rc.local.j2
#!/bin/bash
/opt/solr-{{ solr_version }}/bin/solr start -p {{ solr_port }} -force
exit 0
