# evoacme/files/certbot.cron
#!/bin/bash

ls /etc/letsencrypt/*.crt 2>/dev/null |sed 's/.crt//'|while read vhost;
do 
	evoacme $(basename $vhost)>/dev/null
done
